# src/lib/systemd/system/minecraft-monitoring.timer
[Unit]
Description=Minecraft Server Monitoring

[Timer]
OnActiveSec=30s
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
